# III-ISNR-IB-BRS/publisher_Output file/crontab
# Run every 2 seconds

*/2 * * * * * * python fetchdata.py
